# deploy/systemd/wukong-agent.service
# wukong 探针 systemd 服务
[Unit]
Description=wukong Agent - Server Monitor
Documentation=https://github.com/wukong/monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/wukong/agent/wukong-agent --config /opt/wukong/agent/agent.conf
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/opt/wukong/agent
Restart=always
RestartSec=5

# 探针需要读取系统指标，但不需要高权限
# 若需 ICMP Ping，需 CAP_NET_RAW
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=full
PrivateTmp=yes

# 文件描述符
LimitNOFILE=65536

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wukong-agent

[Install]
WantedBy=multi-user.target